# packages/rclone/.config/systemd/user/obsidian-sync.service
[Unit]
Description=Sync Obsidian Vault with OneDrive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone bisync /home/james/Documents/Obsidian onedrive:Documents/Obsidian --resync --resync-mode newer --create-empty-src-dirs --progress --verbose --transfers 10 --checkers 5

[Install]
WantedBy=default.target
