# packages/rclone/.config/systemd/user/obsidian-sync.service
[Unit]
Description=Sync Obsidian Vault with OneDrive
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/james/.config/systemd/user/obsidian-sync.sh

[Install]
WantedBy=default.target
