# packages/rclone/.config/systemd/user/obsidian-sync.timer
[Unit]
Description=Periodically sync Obsidian Vault with OneDrive

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
